# .config/systemd/user/check-battery-user.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/diego/.config/systemd/user/check-battery.sh
