# .config/systemd/user/org-backup.service
[Unit]
Description=Auto Git Commit and Push for ~/org

[Service]
Type=oneshot
ExecStart=%h/.local/.scripts/automatic/orgBackup.sh
